# deploy/systemd/sqtoimage-backend.service
[Unit]
Description=sqtoimage backend (FastAPI/Uvicorn)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=sqtoimage
Group=sqtoimage
WorkingDirectory=/opt/sqtoimage/backend
EnvironmentFile=/opt/sqtoimage/backend/.env
ExecStart=/opt/sqtoimage/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 9000 --workers 2
Restart=always
RestartSec=3
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/sqtoimage/backend/storage

[Install]
WantedBy=multi-user.target

